# autostart/flashfocus.desktop
[Desktop Entry]
Name=Flash Focus
Comment=Flash window on focus gain
Exec=flashfocus
Terminal=false
Type=Application
