# autostart/mpd.desktop
[Desktop Entry]
Name=Music Player Daemon
Exec=mpd
Terminal=false
Type=Application
